# deploy/systemd/indexkit-update.service
[Unit]
Description=IndexKit updater (pull latest ZIP and deploy static site)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=indexkit
Group=indexkit
Environment=RELEASE_ZIP_URL=https://example.com/your-latest-release.zip
ExecStart=/usr/local/bin/indexkit-update.sh
