# deploy/systemd/indexkit-update.timer
[Unit]
Description=Run IndexKit updater daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
